# crontab
0 0 * * *  bash -c '(( $(date +\%s) / 86400 \% 2 == 0 )) && php artisan transactions:sum'
